# docs/mesos/cni/docker-compose/init_scripts/etc/rc.local
#!/bin/bash
# Record current boot time. Helps with debugging.
date > /tmp/boottime.txt